# termux/service/pkthere/run
#!/data/data/com.termux/files/usr/bin/sh
set -eu

PREFIX="${PREFIX:-/data/data/com.termux/files/usr}"
TERMUX_HOME="${HOME:-${PREFIX%/usr}/home}"
PKTHERE_BIN="${PREFIX}/bin/pkthere"

if [ ! -x "$PKTHERE_BIN" ] && [ -x "${TERMUX_HOME}/.bin/pkthere" ]; then
    PKTHERE_BIN="${TERMUX_HOME}/.bin/pkthere"
fi

TARGET="$(head -n 1 ./env/TARGET | tr -d '[:space:]')"
PROTO="${BIND_PROTO:-UDP}"
ADDR="${BIND_IP:-0.0.0.0}"
PORT="${BIND_PORT:-21080}"

# Optional: configure bind IP and validate it's not empty or blank
if [ -f ./env/BIND_IP ]; then
    ADDR_OVERRIDE="$(head -n 1 ./env/BIND_IP | tr -d '[:space:]')"
    if [ -n "$ADDR_OVERRIDE" ]; then
        ADDR="$ADDR_OVERRIDE"
    fi
fi

# Allow override of bind port via env/BIND_PORT if set
if [ -f ./env/BIND_PORT ]; then
    PORT_OVERRIDE="$(head -n 1 ./env/BIND_PORT | tr -d '[:space:]')"
    if [ -n "$PORT_OVERRIDE" ]; then
        PORT="$PORT_OVERRIDE"
    fi
fi

exec nice -n -20 "$PKTHERE_BIN" --here "${PROTO}:${ADDR}:${PORT}" --there "${TARGET}" --on-timeout exit --timeout-secs 10 --workers 1 --max-payload 0
