# etc/sv/pinger/run
#!/bin/sh
# runit service "pinger": keepalives & latency logging
# Supports SQLite file or PostgreSQL via Unix socket
# -------------------------------------------------------------
# Location: /etc/sv/pinger/run
# Service account: pinger
#
# Prerequisites (run once as root):
# Common:
#   install -d -o pinger -g pinger /var/db/pinger
#   install -o root -g pinger -m 640 etc/pinger.env <OS_NATIVE_PINGER_ENV>
#   # If using PostgreSQL, also set PGHOST, PGPORT, PGDATABASE, PGUSER in the env file.
#
# FreeBSD:
#   pw groupadd -g "$gid" -n "$name"
#   pw useradd -u "$uid" -g "$gid" -d "$home" -w no -s /usr/sbin/nologin -n "$name"
#   # Example:
#   #   name=pinger gid=973 uid=973 home=/nonexistent
#   #   pw groupadd -g "$gid" -n "$name" || true
#   #   pw useradd -u "$uid" -g "$gid" -d "$home" -w no -s /usr/sbin/nologin -n "$name" || true
#   install -d -o root -g wheel -m 755 /usr/local/etc
#   install -o root -g pinger -m 640 etc/pinger.env /usr/local/etc/pinger.env
#
# Linux:
#   groupadd -g "$gid" -r "$name"
#   useradd -u "$uid" -g "$gid" -d "$home" -M -N -r -s /usr/sbin/nologin "$name"
#   # Example:
#   #   name=pinger gid=973 uid=973 home=/nonexistent
#   #   groupadd -g "$gid" -r "$name" || true
#   #   useradd -u "$uid" -g "$gid" -d "$home" -M -N -r -s /usr/sbin/nologin "$name" || true
#   install -d -o root -g root -m 755 /etc
#   install -o root -g pinger -m 640 etc/pinger.env /etc/pinger.env
# -------------------------------------------------------------
# Enable service: ln -s /etc/sv/pinger /var/service/
set -eu

SCRIPT_NAME=pinger
DATE_EPOCH_MODE=""
OS_NAME=""

has_cmd() {
    command -v "$1" >/dev/null 2>&1
}

log() {
    msg="$1"
    if has_cmd logger; then
        logger -t "$SCRIPT_NAME" -p daemon.notice "$msg"
    else
        printf '%s\n' "$msg" >&2
    fi
}

log_warn() {
    msg="$1"
    if has_cmd logger; then
        logger -t "$SCRIPT_NAME" -p daemon.warning "$msg"
    else
        printf '%s\n' "$msg" >&2
    fi
}

log_err() {
    msg="$1"
    if has_cmd logger; then
        logger -t "$SCRIPT_NAME" -p daemon.err "$msg"
    else
        printf '%s\n' "$msg" >&2
    fi
}

log_sql_messages() {
    log_fn="$1"
    operation="$2"
    stderr_file="$3"
    saw_output=false

    while IFS= read -r line || [ -n "$line" ]; do
        saw_output=true
        "$log_fn" "[sql:$operation] $line"
    done <"$stderr_file"

    if [ "$saw_output" != true ]; then
        "$log_fn" "[sql:$operation] command failed"
    fi
}

run_sql_command() {
    sql="$1"
    operation="$2"
    log_fn="$3"
    stderr_file=$(umask 077 && mktemp "/tmp/${SCRIPT_NAME}.sql.XXXXXX") || {
        "$log_fn" "[sql:$operation] failed to create temp file for SQL stderr capture"
        return 1
    }

    if output=$(printf '%s\n' "$sql" | eval "$CLIENT_CMD" 2>"$stderr_file"); then
        rm -f -- "$stderr_file"
        printf '%s' "$output"
        return 0
    fi

    log_sql_messages "$log_fn" "$operation" "$stderr_file"
    rm -f -- "$stderr_file"
    return 1
}

format_sql_timestamp() {
    tick="$1"
    case "$DATE_EPOCH_MODE" in
    bsd)
        date -u -r "$tick" '+%Y-%m-%d %H:%M:%S+00:00'
        ;;
    gnu)
        date -u -d "@$tick" '+%Y-%m-%d %H:%M:%S+00:00'
        ;;
    *)
        log_err "unsupported DATE_EPOCH_MODE='$DATE_EPOCH_MODE'"
        exit 1
        ;;
    esac
}

detect_os_name() {
    OS_NAME=$(uname -s 2>/dev/null || printf 'unknown')
}

detect_date_epoch_mode() {
    probe_tick=0
    if date -u -r "$probe_tick" '+%Y-%m-%d %H:%M:%S' >/dev/null 2>&1; then
        DATE_EPOCH_MODE=bsd
        return 0
    fi
    if date -u -d "@$probe_tick" '+%Y-%m-%d %H:%M:%S' >/dev/null 2>&1; then
        DATE_EPOCH_MODE=gnu
        return 0
    fi
    log_err "unable to determine supported epoch date format"
    exit 1
}

is_number() {
    value="$1"
    case "$value" in
    '' | - | . | -. | *[!0-9.-]* | *.*.* | --* | *-*)
        case "$value" in
        -[0-9]* | -[0-9]*.[0-9]* | [0-9]* | [0-9]*.[0-9]*)
            return 0
            ;;
        *)
            return 1
            ;;
        esac
        ;;
    *)
        return 0
        ;;
    esac
}

sql_numeric_or_null() {
    value="$1"
    if [ -z "$value" ]; then
        printf 'NULL'
        return 0
    fi
    if is_number "$value"; then
        printf '%s' "$value"
        return 0
    fi
    return 1
}

resolve_service_config_path() {
    case "$OS_NAME" in
    FreeBSD)
        printf '/usr/local/etc/pinger.env'
        ;;
    Linux)
        printf '/etc/pinger.env'
        ;;
    *)
        log_err "unsupported OS for config path resolution: '$OS_NAME'"
        exit 1
        ;;
    esac
}

detect_os_name
CONFIG_FILE=$(resolve_service_config_path)
if [ ! -f "$CONFIG_FILE" ]; then
    log_err "missing pinger config file: $CONFIG_FILE"
    exit 1
fi
# shellcheck disable=SC1090
. "$CONFIG_FILE"

# Defaults
: "${SERVICE_USER:=pinger}"
: "${RING_BUFFER_DAYS:=30}"
: "${TARGET:=2001:470:20::2}"
: "${DB_TYPE:=sqlite}"
: "${DB_FILE:=/var/db/pinger/pings.db}"
: "${PGHOST:=/tmp}"
: "${PGPORT:=5432}"
: "${PGDATABASE:=pinger_db}"
: "${PGUSER:=pinger}"

if [ "$(id -un)" != "$SERVICE_USER" ]; then
    log "Dropping privileges to '$SERVICE_USER'"
    exec chpst -u "$SERVICE_USER" /bin/sh "$0" "$@"
fi

RING_BUFFER_SIZE=$((RING_BUFFER_DAYS * 86400))

INIT_SCHEMA_TMPL=$(
    cat <<'EOS'
CREATE TABLE IF NOT EXISTS %TABLE% (
    slot             INTEGER PRIMARY KEY,
    ts               %TS_TYPE% NOT NULL DEFAULT CURRENT_TIMESTAMP UNIQUE,
    latency_ms       REAL NOT NULL
);
EOS
)

ENSURE_INDEXES_TMPL=$(
    cat <<'EOS'
CREATE UNIQUE INDEX IF NOT EXISTS %TS_INDEX% ON %TABLE%(ts);
EOS
)

UPSERT_TMPL=$(
    cat <<'EOS'
INSERT INTO %TABLE%(slot, ts, latency_ms) VALUES(
    %CURRENT_SLOT%, '%CURRENT_TS%', %RTT%
) ON CONFLICT(slot) DO UPDATE SET
    ts               = '%CURRENT_TS%',
    latency_ms       = EXCLUDED.latency_ms;
EOS
)

LAST_SLOT_TMPL="SELECT slot FROM %TABLE% ORDER BY ts DESC LIMIT 1;"

case "$DB_TYPE" in
sqlite)
    TABLE_NAME="pings"
    TS_TYPE="DATETIME"
    TS_INDEX_NAME="pings_ts_idx"
    CLIENT_CMD="sqlite3 '$DB_FILE'"
    ;;
postgres)
    TABLE_NAME="pinger.pings"
    TS_TYPE="TIMESTAMPTZ"
    TS_INDEX_NAME="pings_ts_idx"
    CLIENT_CMD="psql -qAt -h '$PGHOST' -p '$PGPORT' -U '$PGUSER' -d '$PGDATABASE'"
    ;;
*)
    log_err "unsupported DB_TYPE='$DB_TYPE'"
    exit 1
    ;;
esac

render_table_sql() {
    printf '%s\n' "$1" | sed "s/%TABLE%/$TABLE_NAME/g"
}

render_schema_sql() {
    printf '%s\n' "$INIT_SCHEMA_TMPL" |
        sed "s/%TABLE%/$TABLE_NAME/g; s/%TS_TYPE%/$TS_TYPE/g"
}

render_index_sql() {
    printf '%s\n' "$ENSURE_INDEXES_TMPL" |
        sed "s/%TABLE%/$TABLE_NAME/g; s/%TS_INDEX%/$TS_INDEX_NAME/g"
}

init_sql=$(render_schema_sql)
run_sql_command "$init_sql" "init-schema" log_err >/dev/null

index_sql=$(render_index_sql)
run_sql_command "$index_sql" "ensure-indexes" log_err >/dev/null

log "Initialized pinger target='$TARGET' db_type='$DB_TYPE' ring_buffer_days='$RING_BUFFER_DAYS' ring_buffer_size='$RING_BUFFER_SIZE'"

LAST_SLOT_SQL=$(render_table_sql "$LAST_SLOT_TMPL")
LAST_SLOT=$(run_sql_command "$LAST_SLOT_SQL" "select-last-slot" log_err)

case "$LAST_SLOT" in
'' | *[!0-9]*) CURRENT_SLOT=0 ;;
*) CURRENT_SLOT=$(((LAST_SLOT + 1) % RING_BUFFER_SIZE)) ;;
esac

detect_date_epoch_mode
next_tick=$(($(date -u +%s) + 1))
target_reachable=1

while :; do
    while :; do
        now=$(date -u +%s)
        [ "$now" -ge "$next_tick" ] && break
        sleep 0.1
    done
    sample_tick=$next_tick
    next_tick=$((now + 1))
    CURRENT_TS=$(format_sql_timestamp "$sample_tick")

    RTT=$(ping6 -c 1 -W 1 -q "$TARGET" 2>/dev/null | awk -F/ '/[Mm]in\/avg\/max/ {print $5; exit}')

    AVG_CPU_FREQ_MHZ_SQL=$(sql_numeric_or_null "$AVG_CPU_FREQ_MHZ") || {
        log_warn "failed to render avg cpu freq for SQL"
        continue
    }

    UPSERT_SQL=$(printf '%s\n' "$UPSERT_TMPL" |
        sed "s/%TABLE%/$TABLE_NAME/g; s/%CURRENT_SLOT%/$CURRENT_SLOT/; s/%CURRENT_TS%/$CURRENT_TS/; s/%RTT%/$RTT/")
    if ! run_sql_command "$UPSERT_SQL" "write-sample" log_err >/dev/null; then
        continue
    fi
    CURRENT_SLOT=$(((CURRENT_SLOT + 1) % RING_BUFFER_SIZE))
done
