# etc/sv/pinger/run
#!/bin/sh
# runit service "pinger": keepalives & latency logging
# Supports SQLite file or PostgreSQL via Unix socket
# -------------------------------------------------------------
# Location: /etc/sv/pinger/run
# Service account: pinger
#
# Prerequisites (run once as root):
#   pw groupadd -n pinger || true
#   pw useradd -n pinger -g pinger -c "pinger service" -s /usr/sbin/nologin -d /nonexistent -w no || true
#   install -d -o pinger -g pinger /var/db/pinger
#   install -d -o root -g pinger -m 755 /etc/sv/pinger/env
#   install -o root -g pinger -m 640 /dev/null /etc/sv/pinger/env/.env.pinger
#   cat > /etc/sv/pinger/env/.env.pinger <<'EOF'
#   SERVICE_USER=pinger
#   RING_BUFFER_DAYS=30
#   TARGET=2001:470:20::2
#   DB_TYPE=sqlite
#   EOF
#
# If using PostgreSQL, also set PGHOST, PGPORT, PGDATABASE, PGUSER in the .env file.
# -------------------------------------------------------------
# Enable service: ln -s /etc/sv/pinger /var/service/
set -eu

SCRIPT_NAME=pinger
DATE_EPOCH_MODE=""

# Helpers
has_cmd() {
    command -v "$1" >/dev/null 2>&1
}

log() {
    msg="$1"
    if has_cmd logger; then
        logger -t "$SCRIPT_NAME" -p daemon.notice "$msg"
    else
        printf '%s\n' "$msg" >&2
    fi
}

log_err() {
    msg="$1"
    if has_cmd logger; then
        logger -t "$SCRIPT_NAME" -p daemon.err "$msg"
    else
        printf '%s\n' "$msg" >&2
    fi
}

log_sql_messages() {
    log_fn="$1"
    operation="$2"
    stderr_file="$3"
    saw_output=false

    while IFS= read -r line || [ -n "$line" ]; do
        saw_output=true
        "$log_fn" "[sql:$operation] $line"
    done <"$stderr_file"

    if [ "$saw_output" != true ]; then
        "$log_fn" "[sql:$operation] command failed"
    fi
}

run_sql_command() {
    sql="$1"
    operation="$2"
    log_fn="$3"
    stderr_file=$(umask 077 && mktemp "/tmp/${SCRIPT_NAME}.sql.XXXXXX") || {
        "$log_fn" "[sql:$operation] failed to create temp file for SQL stderr capture"
        return 1
    }

    if output=$(printf '%s\n' "$sql" | eval "$CLIENT_CMD" 2>"$stderr_file"); then
        rm -f -- "$stderr_file"
        printf '%s' "$output"
        return 0
    fi

    log_sql_messages "$log_fn" "$operation" "$stderr_file"
    rm -f -- "$stderr_file"
    return 1
}

format_sql_timestamp() {
    tick="$1"
    case "$DATE_EPOCH_MODE" in
    bsd)
        date -u -r "$tick" '+%Y-%m-%d %H:%M:%S+00:00'
        ;;
    gnu)
        date -u -d "@$tick" '+%Y-%m-%d %H:%M:%S+00:00'
        ;;
    *)
        log_err "unsupported DATE_EPOCH_MODE='$DATE_EPOCH_MODE'"
        exit 1
        ;;
    esac
}

detect_date_epoch_mode() {
    probe_tick=0
    if date -u -r "$probe_tick" '+%Y-%m-%d %H:%M:%S' >/dev/null 2>&1; then
        DATE_EPOCH_MODE=bsd
        return 0
    fi
    if date -u -d "@$probe_tick" '+%Y-%m-%d %H:%M:%S' >/dev/null 2>&1; then
        DATE_EPOCH_MODE=gnu
        return 0
    fi
    log_err "unable to determine supported epoch date format"
    exit 1
}

# Environment
ENV_DIR="$(dirname "$0")/env"
# shellcheck disable=SC1091
[ -f "$ENV_DIR/.env.pinger" ] && . "$ENV_DIR/.env.pinger"

# Defaults
: "${SERVICE_USER:=pinger}"
: "${RING_BUFFER_DAYS:=30}"
: "${TARGET:=2001:470:20::2}"
: "${DB_TYPE:=sqlite}"
: "${DB_FILE:=/var/db/pinger/pings.db}"
: "${PGHOST:=/tmp}"
: "${PGPORT:=5432}"
: "${PGDATABASE:=pinger_db}"
: "${PGUSER:=pinger}"

# Startup validation and privilege drop
if [ "$(id -un)" != "$SERVICE_USER" ]; then
    log "Dropping privileges to '$SERVICE_USER'"
    exec chpst -u "$SERVICE_USER" /bin/sh "$0" "$@"
fi

RING_BUFFER_SIZE=$((RING_BUFFER_DAYS * 86400))

# SQL setup
INIT_SCHEMA_TMPL=$(
    cat <<'EOS'
-- SQLite or PostgreSQL schema initialization
CREATE TABLE IF NOT EXISTS %TABLE% (
    slot       INTEGER PRIMARY KEY,
    ts         %TSTYPE% NOT NULL DEFAULT CURRENT_TIMESTAMP UNIQUE,
    latency_ms REAL   NOT NULL
);
EOS
)

UPSERT_TMPL=$(
    cat <<'EOS'
-- SQLite or PostgreSQL upsert
INSERT INTO %TABLE%(slot, ts, latency_ms) VALUES(
    %CURRENT_SLOT%, '%CURRENT_TS%', %RTT%
) ON CONFLICT(slot) DO UPDATE SET
    ts         = '%CURRENT_TS%',
    latency_ms = EXCLUDED.latency_ms;
EOS
)

LAST_SLOT_TMPL="SELECT slot FROM %TABLE% ORDER BY ts DESC LIMIT 1;"

# Determine client command and schema substitutions
case "$DB_TYPE" in
sqlite)
    SED_RENDER='s/%TABLE%/pings/g; s/%TSTYPE%/DATETIME/g'
    CLIENT_CMD="sqlite3 '$DB_FILE'"
    ;;
postgres)
    SED_RENDER='s/%TABLE%/pinger.pings/g; s/%TSTYPE%/TIMESTAMPTZ/g'
    CLIENT_CMD="psql -qAt -h '$PGHOST' -p '$PGPORT' -U '$PGUSER' -d '$PGDATABASE'"
    ;;
*)
    log_err "unsupported DB_TYPE='$DB_TYPE'"
    exit 1
    ;;
esac

# Initialize schema
INIT_SCHEMA_SQL=$(echo "$INIT_SCHEMA_TMPL" | sed "$SED_RENDER")
run_sql_command "$INIT_SCHEMA_SQL" \
    "init-schema" \
    log_err >/dev/null
log "Initialized pinger target='$TARGET' db_type='$DB_TYPE' ring_buffer_days='$RING_BUFFER_DAYS' ring_buffer_size='$RING_BUFFER_SIZE'"

# Determine starting slot
LAST_SLOT_SQL=$(echo "$LAST_SLOT_TMPL" | sed "$SED_RENDER")
LAST_SLOT=$(run_sql_command "$LAST_SLOT_SQL" "select-last-slot" log_err)

case "$LAST_SLOT" in
'' | *[!0-9]*) CURRENT_SLOT=0 ;;
*) CURRENT_SLOT=$(((LAST_SLOT + 1) % RING_BUFFER_SIZE)) ;;
esac

# Sampling setup
detect_date_epoch_mode
target_reachable=1
next_tick=$(($(date -u +%s) + 1))

# Main loop: keep a strict 1 Hz cadence aligned to wall-clock seconds and skip
# missed intervals instead of catching up within the same second.
while :; do
    while :; do
        now=$(date -u +%s)
        [ "$now" -ge "$next_tick" ] && break
        sleep 0.1
    done
    sample_tick=$next_tick
    next_tick=$((now + 1))
    CURRENT_TS=$(format_sql_timestamp "$sample_tick")

    # One probe with a 1-second timeout (BSD ping6 uses seconds for -W)
    RTT=$(ping6 -c 1 -W 1 -q "$TARGET" 2>/dev/null | awk -F/ '/[Mm]in\/avg\/max/ {print $5; exit}')

    if [ -n "$RTT" ]; then
        if [ "$target_reachable" -ne 1 ]; then
            log "target '$TARGET' recovered"
            target_reachable=1
        fi
        UPSERT_SQL=$(echo "$UPSERT_TMPL" |
            sed "$SED_RENDER; s/%CURRENT_SLOT%/$CURRENT_SLOT/; s/%CURRENT_TS%/$CURRENT_TS/; s/%RTT%/$RTT/")
        if ! run_sql_command "$UPSERT_SQL" "write-sample" log_err >/dev/null; then
            continue
        fi
        CURRENT_SLOT=$(((CURRENT_SLOT + 1) % RING_BUFFER_SIZE))
    else
        if [ "$target_reachable" -eq 1 ]; then
            log_err "target '$TARGET' became unreachable"
            target_reachable=0
        fi
    fi
done
