# etc/sv/pinger/run
#!/bin/sh
# runit service "pinger": keepalives & latency logging
# -------------------------------------------------------------
# Location: /etc/sv/pinger/run
# Service account: pinger
#
# Prerequisites (run once as root):
#   pw useradd -n pinger -c "pinger service" -s /usr/sbin/nologin -d /nonexistent -w no
#   mkdir -p /var/db/pinger && chown -R pinger:pinger /var/db/pinger
#   mkdir -p /etc/sv/pinger/env && chmod -R 755 /etc/sv/pinger && \
#     echo "RING_BUFFER_DAYS=7" > /etc/sv/pinger/env/RING_BUFFER_DAYS && \
#     echo "TARGET=2001:470:20::2" > /etc/sv/pinger/env/TARGET
#
# Enable service: ln -s /etc/sv/pinger /var/service/

# Exit on unset variables or errors for reliability and safety
set -eu

# Drop privileges to pinger if not already running as pinger
current_user=$(id -un)
if [ "$current_user" != "pinger" ]; then
    exec chpst -u pinger /bin/sh "$0" "$@"
fi

# Load configuration
: "${RING_BUFFER_DAYS:=$(cat /etc/sv/pinger/env/RING_BUFFER_DAYS 2>/dev/null || echo 7)}"
: "${TARGET:=$(cat /etc/sv/pinger/env/TARGET 2>/dev/null || echo 2001:470:20::2)}"
RING_BUFFER_SIZE=$((RING_BUFFER_DAYS * 86400))
DB="/var/db/pinger/pings.db"

# Sleep for a startup delay to keep ts unique
sleep 1

# Initialize database and schema
sqlite3 "$DB" <<EOF
PRAGMA journal_mode = WAL;
PRAGMA synchronous = NORMAL;
CREATE TABLE IF NOT EXISTS pings (
    slot       INTEGER PRIMARY KEY,
    ts         DATETIME DEFAULT CURRENT_TIMESTAMP UNIQUE,
    latency_ms REAL    NOT NULL
);
CREATE INDEX IF NOT EXISTS idx_pings_ts ON pings(ts);
EOF

# Determine starting slot
last_slot=$(sqlite3 "$DB" 'SELECT slot FROM pings ORDER BY ts DESC LIMIT 1;')
case "$last_slot" in
'' | *[!0-9]*) slot=0 ;;                           # start at zero if empty or invalid
*) slot=$(((last_slot + 1) % RING_BUFFER_SIZE)) ;; # wrap
esac

# Main loop: ping, log, advance
while :; do
    # Ping with native timeout; parse average RTT from min/avg/max line
    rtt=$(ping6 -c 1 -W 1 -q "$TARGET" 2>/dev/null | awk -F/ '/[Mm]in\/avg\/max/ {print $5}')

    if [ -n "$rtt" ]; then
        if sqlite3 "$DB" \
            "INSERT INTO pings (slot, latency_ms) VALUES ($slot, $rtt) \
             ON CONFLICT(slot) DO UPDATE SET ts = CURRENT_TIMESTAMP, latency_ms = excluded.latency_ms;"; then
            slot=$(((slot + 1) % RING_BUFFER_SIZE))
        else
            echo "[pinger] ERROR: Could not insert at slot $slot (possible duplicate timestamp)" >&2
        fi
    fi

    sleep 1
done
