# etc/sv/vpn-runner/run
#!/bin/sh
set -eu

# shellcheck disable=SC2034
SCRIPT_NAME=vpn-runner
SCRIPT_DIR=$(CDPATH='' cd -- "$(dirname -- "$0")" && pwd)
ENV_CONF=/usr/local/etc/.gateway-env.conf
STARTUP_DELAY=10
DHCLIENT_SERVICE=dhclient
OPENVPN_PID_FILE=/var/run/openvpn.pid
OPENVPN_SYSLOG_TAG=openvpn
PROXY_PROBE_TIMEOUT=5
PROXY_PROBE_PORT=5201

# shellcheck disable=SC1091
. "$SCRIPT_DIR/common.sh"

# Delay startup to ensure network is settled.
sleep "$STARTUP_DELAY"

# Ensure the script is run as root.
if [ "$(id -u)" -ne 0 ]; then
    log_err "Error: This script must be run as root. Exiting."
    exit 1
fi

# Load environment variables from /usr/local/etc/.gateway-env.conf config file.
ensure_absolute_path "Environment config path" "$ENV_CONF"
# shellcheck disable=SC1090
. "$ENV_CONF"

ensure_absolute_path "OpenVPN template path" "$OVPN_TEMPLATE"
ensure_absolute_path "OpenVPN config path" "$OVPN_FILE"
ensure_absolute_path "OpenVPN auth path" "$OVPN_AUTH"
[ -z "${OVPN_PROXY_AUTH:-}" ] || ensure_absolute_path "OpenVPN proxy auth path" "$OVPN_PROXY_AUTH"
ensure_safe_interface_list "External interface list" "$EXT_IF_LIST"

# Configuration: timeouts, intervals, and retry settings.
OVPN_PROXY_PORT=$(tail -n 1 "${OVPN_PROXY_AUTH:-}" 2>/dev/null || echo 0)
ensure_nonempty_decimal "OpenVPN proxy port" "$OVPN_PROXY_PORT"
GATEWAY_TIMEOUT=10 # Total seconds to wait for default gateway(s).
INTERVAL=1         # Polling interval in seconds.

# Get current default gateway interface.
EXT_IF_CUR=$(netstat -rn -f inet 2>/dev/null | awk '$1 == "default" {print $4; exit}')
log "Current default interface='${EXT_IF_CUR:-none}' candidates='${EXT_IF_LIST}'"

# Build a newline-separated list of candidates, trying the current default last.
CANDIDATE_LINES_OTHER=$(echo "$EXT_IF_LIST" | tr ',' '\n' | grep -Fvx "${EXT_IF_CUR:-}")
CANDIDATE_LINES_CUR=$(echo "$EXT_IF_LIST" | tr ',' '\n' | grep -Fx "${EXT_IF_CUR:-}" || true)
CANDIDATE_LINES=$(printf '%s\n%s\n' "$CANDIDATE_LINES_OTHER" "$CANDIDATE_LINES_CUR" | sed '/./!d')

if [ -z "$CANDIDATE_LINES" ]; then
    log_err "Error: No candidate interfaces available."
    exit 1
fi

# Loop until we have a working gateway or run out of candidates.
EXT_IF_NEW=""
EXT_IF_AVAILABLE=0
DEFAULT_GW_IPV4=""
while [ -n "$CANDIDATE_LINES" ]; do
    # Pick a random index from all but the last candidate, reserving the last entry
    # so the current default interface is retried only after every other candidate fails.
    NUM_CANDIDATES=$(echo "$CANDIDATE_LINES" | wc -l | tr -d ' ')
    if [ "$NUM_CANDIDATES" -gt 1 ]; then
        RAND=$(od -An -N2 -tu2 /dev/urandom | tr -d ' ')
        INDEX=$((RAND % (NUM_CANDIDATES - 1) + 1))
    else
        INDEX=1
    fi

    # Grab exactly that line.
    EXT_IF_NEW=$(echo "$CANDIDATE_LINES" | sed -n "${INDEX}p")
    ensure_safe_interface_name "Candidate interface" "$EXT_IF_NEW"
    log "Trying candidate interface '$EXT_IF_NEW' index=$INDEX"

    remove_current() {
        CANDIDATE_LINES=$(echo "$CANDIDATE_LINES" | sed "${INDEX}d")
    }

    # Verify the specified interface exists.
    EXT_IF_AVAILABLE=$([ -n "$EXT_IF_NEW" ] && ifconfig "$EXT_IF_NEW" >/dev/null 2>&1 && echo 1 || echo 0)
    if [ "$EXT_IF_AVAILABLE" -ne 1 ]; then
        log_err "Interface '$EXT_IF_NEW' is not available; trying another interface."
        remove_current
        continue
    fi

    # Delete existing default routes for IPv4 and IPv6.
    log "Clearing existing default routes before DHCP on '$EXT_IF_NEW'"
    route -n delete -inet default || true
    route -n delete -inet6 default || true

    # Restart the DHCP client service on the interface.
    log "Restarting dhclient on '$EXT_IF_NEW'"
    if ! timeout -f "$GATEWAY_TIMEOUT" service "$DHCLIENT_SERVICE" restart "$EXT_IF_NEW"; then
        log_err "DHCP restart failed on $EXT_IF_NEW; trying another interface."
        remove_current
        continue
    fi

    # Poll for default gateway on this interface.
    SECONDS_WAITED=0

    while [ "$SECONDS_WAITED" -lt "$GATEWAY_TIMEOUT" ]; do
        DEFAULT_GW_IPV4=$(netstat -rn -f inet 2>/dev/null | awk -v iface="$EXT_IF_NEW" '$1 == "default" && $4 == iface {print $2; exit}')
        if [ -z "$DEFAULT_GW_IPV4" ]; then
            sleep "$INTERVAL"
            SECONDS_WAITED=$((SECONDS_WAITED + INTERVAL))
            continue
        fi
        break
    done

    if [ -z "$DEFAULT_GW_IPV4" ]; then
        log_err "Failed to obtain DHCP default gateway on $EXT_IF_NEW; trying another interface."
        remove_current
        continue
    fi

    log "Selected interface '$EXT_IF_NEW' with default gateway '$DEFAULT_GW_IPV4'"

    # Break out of the loop now that the interface is up.
    break
done

if [ -z "$DEFAULT_GW_IPV4" ]; then
    log_err "Error: Failed to discover a default gateway after trying all interfaces."
    exit 1
fi

if ! timeout "$PROXY_PROBE_TIMEOUT" nc -z -4 "$DEFAULT_GW_IPV4" "$PROXY_PROBE_PORT" >/dev/null 2>&1; then
    # Connect directly to the server.
    OVPN_PROXY_PORT=0
    log "Proxy port probe failed for gateway '$DEFAULT_GW_IPV4'; using direct OpenVPN connection"
else
    log "Proxy port probe succeeded for gateway '$DEFAULT_GW_IPV4' port='$OVPN_PROXY_PORT'"
fi

# Update the OpenVPN configuration.
if [ "$EXT_IF_AVAILABLE" -eq 1 ] && [ -f "${OVPN_PROXY_AUTH:-}" ] &&
    [ -n "$DEFAULT_GW_IPV4" ] && [ "$OVPN_PROXY_PORT" -gt 0 ]; then
    if ! write_file_atomically "$OVPN_FILE" \
        sed -e "s|__AUTH_FILE__|$OVPN_AUTH|" \
        -e "s|__PROXY_STRING__|${DEFAULT_GW_IPV4} ${OVPN_PROXY_PORT}|" \
        -e "s|__PROXY_AUTH_FILE__|$OVPN_PROXY_AUTH|" \
        "$OVPN_TEMPLATE"; then
        log_err "Failed to write OpenVPN config '$OVPN_FILE'"
        exit 1
    fi
    log "Wrote OpenVPN config '$OVPN_FILE' using proxy '$DEFAULT_GW_IPV4:$OVPN_PROXY_PORT'"
else
    if ! write_file_atomically "$OVPN_FILE" \
        sed -e "s|__AUTH_FILE__|$OVPN_AUTH|" \
        -e "/__PROXY_STRING__/d" \
        -e "/__PROXY_AUTH_FILE__/d" \
        "$OVPN_TEMPLATE"; then
        log_err "Failed to write OpenVPN config '$OVPN_FILE'"
        exit 1
    fi
    log "Wrote OpenVPN config '$OVPN_FILE' without proxy"
fi

# Start the OpenVPN service with elevated priority to reduce latency.
if command -v rtprio >/dev/null 2>&1; then
    # Realtime class priority 10 (0 highest).
    log "Starting OpenVPN with rtprio priority 10 using config '$OVPN_FILE'"
    exec rtprio 10 openvpn --syslog "$OPENVPN_SYSLOG_TAG" --config "$OVPN_FILE" --writepid "$OPENVPN_PID_FILE"
else
    # Fallback to highest nice if rtprio is unavailable.
    log "Starting OpenVPN with nice -20 using config '$OVPN_FILE'"
    exec nice -n -20 openvpn --syslog "$OPENVPN_SYSLOG_TAG" --config "$OVPN_FILE" --writepid "$OPENVPN_PID_FILE"
fi
