# etc/cron.weekly/pkg-updates
#!/bin/sh
# Update all packages
xbps-install -Syu
flatpak update --noninteractive
