# termux/service/iperf3/run
#!/data/data/com.termux/files/usr/bin/sh
exec nice -n -20 iperf3 -s -1
