# termux/service/iperf3/run
#!/data/data/com.termux/files/usr/bin/sh
printf '%s\n' "[termux:iperf3] starting single-shot iperf3 server" >&2
exec nice -n -20 iperf3 -s -1
