# termux/service/iperf3/run
#!/data/data/com.termux/files/usr/bin/sh
exec iperf3 -s
