# termux/service/3proxy/run
#!/data/data/com.termux/files/usr/bin/sh
exec nice -n -20 ~/3proxy/bin/3proxy ~/.3proxy.cfg
