# etc/sv/iperf3/run
#!/bin/sh
exec chpst -u iperf3 iperf3 -s
